# etc/chantal.service
[Unit]
Description=Kevin CI In-VM agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m chantal
Restart=on-abnormal
RestartForceExitStatus=SIGSEGV SEGV SIGILL SIGKILL

[Install]
WantedBy=multi-user.target
